# pi/netmon.cron.example
# /etc/cron.d/netmon — Sammellauf des netmon-Collectors alle 5 Minuten.
# Startminute bewusst um 2 versetzt, damit er nicht mit netscan.sh
# (läuft auf 0/5/10/...) um die Datenbank konkurriert.
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
2-59/5 * * * * root /usr/local/sbin/netmon_collector.py >> /var/log/netmon/collector.log 2>&1
